# redgreen.service
[Unit]
Description=RedGreen

[Service]
ExecStart=/usr/local/bin/redgreen serve --static /var/www/redgreen

[Install]
WantedBy=multi-user.target
